# agent/DEBIAN/prerm
#!/bin/bash
set -e
systemctl stop astra-monitor.service || true
systemctl disable astra-monitor.service || true
exit 0
