# agent/opt/astra-monitor/astra-monitor.service
[Unit]
Description=Astra File Monitor Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/astra-monitor/agent.py
Restart=always
RestartSec=30
User=root
Group=root

[Install]
WantedBy=multi-user.target
